# development-tripsearchio/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx

[program:mongodb]
command=/usr/bin/mongod --config /etc/mongodb.conf

[program:redis]
command=/usr/bin/redis-server
